# s6-rc.d/mongodb/run
#!/bin/sh
echo "Starting MongoDB..."
exec 2>&1
exec /usr/sbin/mongod --bind_ip 0.0.0.0