# services.d/studybot/run
#!/usr/bin/with-contenv bash
sleep 5 # Allow MongoDB to start
echo "Starting StudyBot..."
exec /usr/bin/npm run start --prefix /app