# s6-rc.d/studybot/run
#!/bin/sh
echo "Starting StudyBot..."
exec /bin/npm run start --prefix /app